# contrib/com.vindu.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- LaunchAgent: cp to ~/Library/LaunchAgents/ and `launchctl load` it so
     vindud starts at login. Adjust the binary path if you installed
     somewhere other than /usr/local/bin. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vindu.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/vindud</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/vindu.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/vindu.log</string>
</dict>
</plist>
